# contrib/gokumail.service
[Unit]
Description=Simple POP3 to IMAP proxy for KUmail

[Service]
ExecStart=/usr/bin/gokumail

[Install]
WantedBy=multi-user.target
